# docker/runtime/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0770
chown=root:picoclaw

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:suprclaw-engine-gateway]
command=su - picoclaw -s /bin/bash -lc "SUPRCLAW_HOME=/home/picoclaw/.suprclaw SUPRCLAW_CONFIG=/home/picoclaw/.suprclaw/config.json SUPRCLAW_ENGINE_HOME=/home/picoclaw/.suprclaw SUPRCLAW_ENGINE_CONFIG=/home/picoclaw/.suprclaw/config.json suprclaw gateway"
autorestart=true
priority=10
stdout_logfile=/var/log/suprclaw-engine-gateway.log
stderr_logfile=/var/log/suprclaw-engine-gateway.err
